# template-storm-supervisord.conf
[program:storm-nimbus]
command=$STORM_INSTALL_DIR/bin/storm nimbus
stopasgroup=true
killasgroup=true
autostart=false		; Because we want the storm-cluster.sh script to start and stop it on demand, not supervisord
autorestart=true	; Autorestart the service if it goes down
user=$STORM_USER

[program:storm-supervisor]
command=$STORM_INSTALL_DIR/bin/storm supervisor
stopasgroup=true
killasgroup=true
autostart=false		; Because we want the storm-cluster.sh script to start and stop it on demand, not supervisord
autorestart=true	; Autorestart the service if it goes down
user=$STORM_USER

[program:storm-logviewer]
command=$STORM_INSTALL_DIR/bin/storm logviewer
stopasgroup=true
killasgroup=true
autostart=false		; Because we want the storm-cluster.sh script to start and stop it on demand, not supervisord
autorestart=true	; Autorestart the service if it goes down
user=$STORM_USER

[program:storm-ui]
command=$STORM_INSTALL_DIR/bin/storm ui
stopasgroup=true
killasgroup=true
autostart=false		; Because we want the storm-cluster.sh script to start and stop it on demand, not supervisord
autorestart=true	; Autorestart the service if it goes down
user=$STORM_USER

